# roles/hl.bind/templates/bind_exporter.service.j2
{{ ansible_managed | comment }}
# Generated by Ansible role {{ ansible_role_name }}

[Unit]
Description=Prometheus
Documentation=https://github.com/prometheus-community/bind_exporter
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User={{ prometheus_user }}
Group={{ prometheus_user }}
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/usr/local/bin/{{ bind_exporter_binary }} \
  --bind.pid-file=/var/run/named/named.pid \
  --bind.timeout=20s \
  --web.listen-address=0.0.0.0:{{ bind_exporter_port }} \
  --web.telemetry-path=/metrics \
  --bind.stats-url=http://127.0.0.1:{{ bind_statistics_port }}/ \
  --bind.stats-groups=server,view,tasks

SyslogIdentifier={{ prometheus_user }}
Restart=always

[Install]
WantedBy=multi-user.target
