# roles/hl.dellhw-exporter/templates/dellhw_exporter.service.j2
{{ ansible_managed | comment }}
# Generated by Ansible role {{ ansible_role_name }}

[Unit]
Description=Prometheus
Documentation=https://github.com/galexrt/dellhw_exporter
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User={{ prometheus_user }}
Group={{ prometheus_user }}
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/usr/local/bin/{{ dellhw_exporter_binary }} \
    --collectors-enabled=chassis,nics,system,memory,processors \
    --collectors-omreport="/opt/dell/srvadmin/bin/omreport" \
    --log-level=INFO \
    --web-listen-address=0.0.0.0:{{ dellhw_exporter_port }} \
    --web-telemetry-path=/metrics

SyslogIdentifier=prometheus
Restart=always

[Install]
WantedBy=multi-user.target
