# roles/hl.haproxy/templates/haproxy_exporter.service.j2
{{ ansible_managed | comment }}
# Generated by Ansible role {{ ansible_role_name }}

[Unit]
Description=Prometheus
Documentation=https://github.com/prometheus/haproxy_exporter
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User={{ prometheus_user }}
Group={{ prometheus_user }}
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/usr/local/bin/{{ haproxy_exporter_binary }} \
  --haproxy.pid-file=/var/run/haproxy.pid \
  --haproxy.timeout=20s \
  --web.listen-address=0.0.0.0:{{ haproxy_exporter_port }} \
  --web.telemetry-path=/metrics \
  '--haproxy.scrape-uri=http://{{ haproxy_stats_user }}:{{ haproxy_stats_password }}@127.0.0.1:{{ haproxy_stats_port }}/stats;csv'

SyslogIdentifier=prometheus
Restart=always

[Install]
WantedBy=multi-user.target
